# docker/python/smart-base.conf
[program:smart-base]
command=python /usr/bin/multi_receive.py
directory=/etc/openvpn
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0