# docker/supervisord/conf/master.conf
[program:master-reveiver]
command=python /home/pi/smart-base/multi_receive.py